# scripts/gilgal-nginx-sync.service
[Unit]
Description=Gilgal Nginx Auto-Sync Service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/gilgal/scripts/nginx-auto-sync.sh
WorkingDirectory=/opt/gilgal
User=azureuser
Group=azureuser

[Install]
WantedBy=multi-user.target